# .config/autostart/rotation.desktop
[Desktop Entry]
Type=Application
Name=rotation
Exec=/home/spdskatr/.xprofile.d/rotate.py2
StartupNotify=false
Terminal=false
